# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:celeryd]
user=api
command=celery -A overlord.celery worker --loglevel=info --beat ;
directory=/var/apps/overlord/overlord ;
autostart=true
autorestart=true
stdout_logfile=/var/apps/overlord/log/out_celeryd.log
stderr_logfile=/var/apps/overlord/log/error_celeryd.log

[program:overlord]
user=api
command=python overlord.py ;
directory=/var/apps/overlord/overlord ;
autostart=true
autorestart=true
stdout_logfile=/var/apps/overlord/log/out_overlord.log
stderr_logfile=/var/apps/overlord/log/error_overlord.log


[program:flower]
user=api
command=flower --port=5555 --basic_auth=tnyu:crushit123 ;
directory=/home/api/overlord/overlord ;
autostart=true
autorestart=true
stdout_logfile=/var/apps/overlord/log/out_celeryd_flower.log
sterr_logfile=/var/apps/overlord/log/error_celeryd_flower.log

[inet_http_server]
port=127.0.0.1:9001   ;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001 ;
